# supervisor_docker.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:gunicorn]
command= gunicorn --bind 0.0.0.0:8080 --workers 3 wsgi:app
directory= /fype03
autostart=true
autorestart=true
